# deploy/ec2/caddy.service
[Unit]
Description=Caddy TLS Reverse Proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/caddy run --config /etc/caddy/Caddyfile
Restart=always
RestartSec=5
Environment=HOME=/root
EnvironmentFile=/etc/opensandbox/caddy.env

[Install]
WantedBy=multi-user.target
